# dbus/tukitd.service
[Unit]
Description=tukitd deamon for transactional update

[Service]
Type=dbus
BusName=org.opensuse.tukit
ExecStart=/usr/sbin/tukitd