# systemd/prepare-nextroot-for-softreboot.service.in
[Unit]
Description=Mount next snapshot for soft-reboot
DefaultDependencies=no
Before=systemd-soft-reboot.service

[Service]
ExecStart=@libexecdir@/prepare-nextroot-for-softreboot
Type=oneshot
RemainAfterExit=yes

[Install]
RequiredBy=soft-reboot.target
